# apps/mlflow/mlflow-server.conf
[program:mlflow-server]
command=mlflow ui --host 0.0.0.0
autostart=true
autorestart=true
stderr_logfile=/var/log/long.err.log
stdout_logfile=/var/log/long.out.log
